# systemd/youtube-webdav-bot.service
[Unit]
Description=YouTube WebDAV Bot
After=network.target

[Service]
Type=simple
User=youtube-bot
WorkingDirectory=/opt/youtube-webdav-bot
Environment="PATH=/opt/youtube-webdav-bot/venv/bin"
ExecStart=/opt/youtube-webdav-bot/venv/bin/python -m bot.main
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/youtube-webdav-bot /var/log/youtube-webdav-bot /tmp/youtube-webdav-bot

[Install]
WantedBy=multi-user.target
